# deploy/com.healthos.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.healthos.server</string>

  <!-- Run uvicorn straight from the venv (no shell activation needed). -->
  <key>ProgramArguments</key>
  <array>
    <string>/Users/node/sandbox/healthos/.venv/bin/uvicorn</string>
    <string>healthos.main:app</string>
    <string>--host</string>
    <string>0.0.0.0</string>
    <string>--port</string>
    <string>8000</string>
  </array>

  <!-- So pydantic finds .env and the app resolves the project root. -->
  <key>WorkingDirectory</key>
  <string>/Users/node/sandbox/healthos</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/node/sandbox/healthos/.venv/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <!-- Start at login and restart on crash = always-on. -->
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <!-- Don't hammer-restart if it's failing fast. -->
  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>StandardOutPath</key>
  <string>/Users/node/sandbox/healthos/logs/server.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/node/sandbox/healthos/logs/server.err.log</string>
</dict>
</plist>
